# backend/debian/postinst
#!/bin/sh
set -e

# Create system users and groups
if ! getent group floxapp >/dev/null; then
    addgroup --system floxapp
fi

if ! id floxsvc >/dev/null 2>&1; then
    adduser --system --ingroup floxapp --no-create-home \
            --disabled-password --shell /bin/false floxsvc
fi

# Create directories
mkdir -p /var/vhosts/flox/{config/environment,generated,logs}
chown -R floxsvc:floxapp /var/vhosts/flox
chmod 2775 /var/vhosts/flox  # SGID for group inheritance

mkdir -p /etc/flox
chown root:floxapp /etc/flox
chmod 0775 /etc/flox

# Set up sudoers
cat > /etc/sudoers.d/flox-backend <<EOF
# Allow floxsvc to run necessary commands
floxsvc ALL=(root) NOPASSWD: /usr/bin/mysql
floxsvc ALL=(root) NOPASSWD: /bin/systemctl reload nginx
floxsvc ALL=(root) NOPASSWD: /bin/ln -s /etc/nginx/sites-available/* /etc/nginx/sites-enabled/
floxsvc ALL=(root) NOPASSWD: /bin/rm /etc/nginx/sites-enabled/*
floxsvc ALL=(root) NOPASSWD: /usr/sbin/nginx -t
floxsvc ALL=(floxapp) NOPASSWD: /var/vhosts/flox/restart.sh
EOF
chmod 0440 /etc/sudoers.d/flox-backend

# Create default config if not exists
if [ ! -f /etc/flox/backend.conf ]; then
    cat > /etc/flox/backend.conf <<EOF
# Flox backend configuration
listen_address: ":8080"
db_admin_path: /etc/flox/mysql-admin.cnf
dns_api_token: ""  # deSEC API token
dns_domain: "flox.click"
template_dir: /var/vhosts/flox/environment/templates
script_dir: /usr/lib/flox/scripts
EOF
    chown root:floxapp /etc/flox/backend.conf
    chmod 0640 /etc/flox/backend.conf
fi

# Create MySQL admin credentials if not exists
if [ ! -f /etc/flox/mysql-admin.cnf ]; then
    cat > /etc/flox/mysql-admin.cnf <<EOF
[client]
user=floxadmin
password=$(openssl rand -base64 24 | tr -dc 'a-zA-Z0-9' | head -c 32)
EOF
    chown root:floxapp /etc/flox/mysql-admin.cnf
    chmod 0640 /etc/flox/mysql-admin.cnf
fi

# Enable and start service
systemctl daemon-reload
systemctl enable flox-backend.service
systemctl start flox-backend.service

echo "Flox backend installed successfully!"
echo "1. Configure MySQL admin user:"
echo "   mysql -u root -p -e \"CREATE USER IF NOT EXISTS 'floxadmin'@'localhost' IDENTIFIED BY 'password'; GRANT ALL ON *.* TO 'floxadmin'@'localhost';\""
echo "2. Update /etc/flox/mysql-admin.cnf with the correct credentials"
