# backend/debian/prerm
#!/bin/sh
set -e

# Stop service before removal
if [ "$1" = "remove" ]; then
    systemctl stop flox-backend.service || true
    systemctl disable flox-backend.service || true
fi
